# crontab
*/5 * * * * /usr/bin/python /path/to/checkPayment.py >> /var/log/checkPayment.log 2>&1